# deploy/systemd/keepa-updater-ca.service
[Unit]
Description=Keepa sheet updater - CA marketplace
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=keepa
WorkingDirectory=/home/keepa/atlas
EnvironmentFile=/home/keepa/atlas/.env
ExecStart=/home/keepa/atlas/.venv/bin/python /home/keepa/atlas/src/main.py update-keepa-sheets --marketplace CA --max-asins 50
StandardOutput=journal
StandardError=journal
